# docker/ai/py36-webssh-cu90/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=service ssh start

[program:gateone]
command=service gateone start